# nas_files/settings/etrayz-settings-init
#!/bin/sh
### BEGIN INIT INFO
# Provides:          etrayz-settings
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:
# Short-Description: Apply EtrayZ settings and startup beep
### END INIT INFO

CONF=/etc/etrayz/settings.conf

case "$1" in
  start)
    [ -f "$CONF" ] || cp /etc/etrayz/defaults.conf "$CONF"
    . "$CONF"

    # Apply all settings
    /usr/local/bin/etrayz-apply all >/dev/null 2>&1

    # Startup beep
    if [ "$BEEP_STARTUP" = "1" ]; then
        /usr/local/bin/etrayz-buzzer "$BEEP_PATTERN" &
    fi

    # Start USB copy button daemon if enabled
    if [ "$USBCOPY_ENABLED" = "1" ]; then
        start-stop-daemon --start --background --make-pidfile \
            --pidfile /var/run/etrayz-usbcopy.pid \
            --exec /usr/local/bin/etrayz-usbcopy 2>/dev/null || true
    fi
    ;;
  stop)
    start-stop-daemon --stop --pidfile /var/run/etrayz-usbcopy.pid 2>/dev/null || true
    ;;
esac
exit 0
